# root/etc/services.d/opentracker/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc opentracker -P $PORT -p $PORT -d /mnt
